# docker/supervisord-onpod.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:redis]
command=/bin/sh -c "exec redis-server --requirepass \"$REDIS_PASSWORD\" --port 6379 --bind 127.0.0.1 --dir /tmp/redis-data --maxmemory 192mb --maxmemory-policy allkeys-lru --save ''"
priority=5
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:redis-rest]
command=node /rest/redis-rest-proxy.mjs
directory=/rest
priority=10
autostart=true
autorestart=true
environment=PORT="8079"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:ais-relay]
command=node /app/scripts/ais-relay.cjs
directory=/app
priority=15
autostart=true
autorestart=true
environment=PORT="3004",NODE_OPTIONS="--max-old-space-size=512"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:kcg-seeders]
command=/bin/sh /app/scripts/kcg-country-seeders.sh
directory=/app
priority=18
autostart=true
autorestart=true
environment=NODE_OPTIONS="--max-old-space-size=256"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worldmonitor-api]
command=node /app/local-api-server.mjs
directory=/app
priority=20
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -c /tmp/nginx.conf -g "daemon off;"
priority=25
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
